# deb-files/debian/postinst
#!/bin/sh
# postinst script for passhport
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see https://www.debian.org/doc/debian-policy/ or
# the debian-policy package

PASSHPORT_HOME='/var/lib/passhport'

case "$1" in
    configure)
	    echo "Creating system user 'passhport'…"
	    adduser --system --group --home "${PASSHPORT_HOME}" --shell /bin/false passhport
	    chown -R passhport:passhport "${PASSHPORT_HOME}" /etc/passhport /var/log/passhport
	    echo "Creating ssh pubkeys (rsa and ecdsa)…"
	    su - passhport -c '/usr/bin/ssh-keygen -t rsa -b 4096 -N "" -f "${PASSHPORT_HOME}/.ssh/id_rsa"'
	    su - passhport -c '/usr/bin/ssh-keygen -t ecdsa -b 521 -N "" -f "${PASSHPORT_HOME}/.ssh/id_ecdsa"'
	    echo "Creating passhport DB…"
	    /var/lib/passhport/python-run-env/bin/python "${PASSHPORT_HOME}/passhport/passhportd/db_create.py"
	    echo "Creating cert private key…"
	    su - passhport -c "openssl genrsa -out "/home/passhport/certs/key.pem" 4096"
	    openssl req -new -key "/home/passhport/certs/key.pem" \
		-config "/usr/share/passhport/openssl.cnf" \
		-out "/home/passhport/certs/cert.pem" \
		-subj "/C=FR/ST=Ile De France/L=Ivry sur Seine/O=LibrIT/OU=DSI/CN=passhport.librit.fr" \
		-x509 \
		-days 365 \
		-sha256 \
		-extensions v3_req
	    systemctl daemon-reload
	    systemctl enable passhportd
	    systemctl start passhportd
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
